# car-panel.service
[Unit]
Description=My Service
After=network.target

[Service]
Type=simple
User=your_username
Group=your_groupname
ExecStart=/path/to/your/program

[Install]
WantedBy=multi-user.target